# crontab/crontab
@daily root curator /elasticsearch/crontab/delete_index.yml --config /elasticsearch/crontab/curator.yml --dry-run
